# torrentwatcher.ini
[supervisord]
nodaemon=true
logfile=/opt/torrentwatcher/log/supervisor.log
loglevel=debug
root=true

[program:openvpn]
autostart=false
directory=/opt/vpn/
command=/usr/sbin/openvpn --config /opt/vpn/user.ovpn --askpass /opt/vpn/key.pass

[program:transmission]
autostart=false
command=/usr/bin/transmission-daemon -f --log-error --logfile /opt/torrentwatcher/log/transmission.log -a '*'

[program:torrentwatcher]
autostart=false
user=watcher
directory=/opt/torrentwatcher
killasgroup=true
command=./torrentwatcher.sh -f /etc/tw.conf

