# torrentwatcher.ini
[supervisord]
nodaemon=true
logfile=/opt/torrentwatcher/log/supervisor.log
loglevel=debug
root=true

[program:openvpn]
autostart=false
directory=/opt/vpn/
command=/usr/sbin/openvpn --config /opt/vpn/user.ovpn --askpass /opt/vpn/key.pass

[program:transmission-daemon]
autostart=false
command=/usr/bin/transmission-daemon -f --log-error --logfile /opt/torrentwatcher/log/transmission.log -a '*'

[program:torrentwatcher]
autostart=true
user=watcher
directory=/opt/torrentwatcher
killasgroup=true
environment=HOME="/opt/",USER="watcher",RCLONE_CONFIG=/etc/rclone.conf
command=/opt/torrentwatcher/torrentwatcher.sh -f /etc/tw.conf

